# bin/home-services.service
[Unit]
Description=Home Services

[Service]
ExecStart=/usr/bin/bash /var/lib/tomcat9/webapps/home-services/bin/run-all.sh
WorkingDirectory=/var/lib/tomcat9/webapps/home-services
StandardOutput=inherit
StandardError=inherit
Restart=no
User=brettonw
RemainAfterExit=true
Type=oneshot

[Install]
WantedBy=multi-user.target
